# postgresql/rootfs/etc/services.d/postgresql/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid postgres /usr/bin/postgres ${POSTGRESQL_OPTIONS}
