# nginx/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/nginx ${NGINX_OPTIONS}
