# mariadb/rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv sh

exec /usr/bin/mysqld --user=mysql --console ${MARIADB_OPTIONS}
